# .config/systemd/user/mpd.socket
[Socket]
ListenStream=7700

[Install]
WantedBy=xlauncher.target
