# .config/systemd/user/xlauncher.target
[Unit]
Description=Xorg launcher
Wants=wm.target

[Install]
Alias=default.target
